# supervisor/rfidreader.conf
[program:RFIDreader]
Command=python3 o.py
directory=/home/pi
autostart=true
autorestart=true
priority=10
startretries=10
user=pi
