# supervisor/robot.conf
[program:Robot]
Command=python3 app.py
directory=/home/pi/programiaki/robot
autostart=true
autorestart=true
priority=10
startretries=10
user=pi
