# supervisor/ws.conf
[program:WebSocketServer]
Command=/home/pi/.nvm/versions/node/v7.10.1/bin/node server.js
directory=/home/pi/programiaki/websockets
autostart=true
autorestart=true
priority=5
startretries=10
user=pi
